# .scripts/com.sw_os.learning-review.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sw_os.learning-review</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>{{VAULT_PATH}}/.scripts/learning-review-prompt.sh</string>
    </array>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>17</integer> <!-- 5 PM -->
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>StandardOutPath</key>
    <string>{{VAULT_PATH}}/.scripts/logs/learning-review-stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>{{VAULT_PATH}}/.scripts/logs/learning-review-stderr.log</string>
    
    <key>WorkingDirectory</key>
    <string>{{VAULT_PATH}}</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
